# terminusa-terminal.service
[Unit]
Description=Terminusa Online Terminal Service
After=network.target postgresql.service redis-server.service
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=terminusa
Group=terminusa
WorkingDirectory=/opt/terminusa
Environment=FLASK_ENV=production
EnvironmentFile=/opt/terminusa/.env

# Start the terminal server
ExecStart=/opt/terminusa/venv/bin/python terminal_server.py

# Restart on failure
Restart=always
RestartSec=5

# Basic security measures
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
