# terminusa-terminal.service
[Unit]
Description=Terminusa Online Terminal Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=terminusa
Group=terminusa
WorkingDirectory=/opt/terminusa
Environment=PYTHONUNBUFFERED=1
Environment=FLASK_ENV=production
EnvironmentFile=/opt/terminusa/.env
ExecStart=/opt/terminusa/venv/bin/python terminal_server.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/terminusa/terminal.log
StandardError=append:/var/log/terminusa/terminal.error.log

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
