# services/terminusa-terminal.service
[Unit]
Description=Terminusa Online Terminal Server
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=terminusa
Group=terminusa
WorkingDirectory=/var/www/terminusa
Environment=PATH=/var/www/terminusa/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/var/www/terminusa
ExecStart=/var/www/terminusa/venv/bin/python terminal_server.py
Restart=always
RestartSec=3
StandardOutput=append:/var/log/terminusa/terminal.log
StandardError=append:/var/log/terminusa/terminal.error.log

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
